# scripts/resourceagent.service
[Unit]
Description=ResourceAgent Monitoring Service
Documentation=https://github.com/your-org/resourceagent
After=network.target

[Service]
Type=simple
User=resourceagent
Group=resourceagent
ExecStart=/opt/resourceagent/resourceagent -config /opt/resourceagent/configs/config.json
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/resourceagent/logs
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
