# scripts/resourceagent.service
[Unit]
Description=ResourceAgent Monitoring Service
After=network.target

[Service]
Type=simple
User=resourceagent
Group=resourceagent
WorkingDirectory=/opt/EEGAgent
ExecStart=/opt/EEGAgent/bin/x86/resourceagent -config /opt/EEGAgent/conf/ResourceAgent/ResourceAgent.json -monitor /opt/EEGAgent/conf/ResourceAgent/Monitor.json -logging /opt/EEGAgent/conf/ResourceAgent/Logging.json
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/EEGAgent/log/ResourceAgent
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
